# main_py.conf
[program:main_py]
command=/bin/bash -c 'cd /path/to/your/project && ./run_main.sh'
directory=/path/to/your/project
user=your_username
autostart=true
autorestart=true
startsecs=5
startretries=3
stopasgroup=true
killasgroup=true
stderr_logfile=/var/log/main_py.err.log
stdout_logfile=/var/log/main_py.out.log
environment=PATH="/path/to/your/project/venv/bin:%(ENV_PATH)s"


# etc/supervisor/conf.d/main_py.conf
# # 1. Install supervisor if not present
# sudo apt update && sudo apt install -y supervisor

# # 2. Create log directory
# sudo mkdir -p /var/log
# sudo touch /var/log/main_py.err.log /var/log/main_py.out.log
# sudo chmod 666 /var/log/main_py.{err,out}.log

# # 3. Copy config and reload
# sudo cp /path/to/your/project/supervisor_main_py.conf /etc/supervisor/conf.d/main_py.conf
# sudo supervisorctl reread
# sudo supervisorctl update

# # 4. Start and check status
# sudo supervisorctl start main_py
# sudo supervisorctl status main_py

# # 5. View logs
# sudo tail -f /var/log/main_py.out.log
# sudo tail -f /var/log/main_py.err.log